# production/ml-nodes/jobanalyzer.cron
# Crontab to run on ml nodes other than ml4.  This just runs sonar every five minutes
# and captures its output.
0-59/5 * * * * $HOME/sonar/sonar.sh
